# ansible/templates/kafka-connect.service
[Unit]
Description=Kafka Service connect

[Service]
LimitNPROC=400000
LimitNOFILE=400000
#Environment="KAFKA_OPTS=-Djava.security.auth.login.config=/opt/kafka/jaas-kafka-server.conf"
Environment="KAFKA_HEAP_OPTS=-Xmx512m -Xms512m"
Type=simple
User=kafka
Group=kafka
ExecStart=/opt/kafka/bin/connect-distributed.sh /opt/kafka/config/connect-distributed.properties
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
