# ansible/templates/kafka-broker.service
[Unit]
Description=Kafka Service broker

[Service]
LimitNPROC=400000
LimitNOFILE=400000
#Environment="KAFKA_OPTS=-Djava.security.auth.login.config=/opt/kafka/jaas-kafka-server.conf"
Environment="KAFKA_HEAP_OPTS=-Xmx512m -Xms512m"
Type=simple
User=kafka
Group=kafka
ExecStart=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/kraft/broker.properties
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
